# tools/docker/frr/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:frr]
command=/usr/lib/frr/frr start 
